# deploy/tg-federal-bot.service
[Unit]
Description=TG Federal Bot
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/tg-federal-bot
EnvironmentFile=-/opt/tg-federal-bot/.env
Environment="JAVA_OPTS=-Xms512m -Xmx1024m"
ExecStart=/usr/bin/java ${JAVA_OPTS} -jar /opt/tg-federal-bot/tg-federal-bot.jar
Restart=always
RestartSec=10
LimitNOFILE=65536

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/tg-federal-bot/logs

[Install]
WantedBy=multi-user.target
